# dot_config/systemd/user/udiskie.service.d/override.conf
[Service]
ExecStart=/usr/bin/udiskie --appindicator
